# linux/shared/cosmic-greeter/cosmic-greeter-wallpaper-sync.path
[Unit]
Description=Watch COSMIC wallpaper config for changes

[Path]
PathChanged=%h/.config/cosmic/com.system76.CosmicBackground
Unit=cosmic-greeter-wallpaper-sync.service

[Install]
WantedBy=default.target
